# deploy/systemd/repociv.target
[Unit]
Description=RepoCiv Imperial Console (Bridge + Frontend)
Documentation=file:%h/.hermes/workspace/repos/repociv/deploy/systemd/README.md
Requires=repociv-bridge.service repociv-frontend.service
After=repociv-bridge.service repociv-frontend.service

[Install]
WantedBy=default.target
