# deploy/systemd/repociv-frontend.service
[Unit]
Description=RepoCiv Frontend (Vite Dev Server)
Documentation=file:%h/.hermes/workspace/repos/repociv/deploy/systemd/README.md
After=network.target repociv-bridge.service
Requires=repociv-bridge.service

[Service]
Type=simple
WorkingDirectory=%h/.hermes/workspace/repos/repociv
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=5
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=development

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=repociv-frontend

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=default.target
