# deploy/systemd/repociv-bridge.service
[Unit]
Description=RepoCiv Bridge Server (Backend API)
Documentation=file:%h/.hermes/workspace/repos/repociv/deploy/systemd/README.md
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.hermes/workspace/repos/repociv
ExecStart=/usr/bin/python3 -m server.bridge
Restart=on-failure
RestartSec=5
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/.hermes/workspace/repos/repociv

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=repociv-bridge

# Security hardening (opcional)
NoNewPrivileges=true

[Install]
WantedBy=default.target
