# deploy/systemd/repociv-backup.service
[Unit]
Description=RepoCiv event store backup
Documentation=file:%h/.hermes/workspace/repos/repociv/deploy/systemd/README.md

[Service]
Type=oneshot
ExecStart=%h/.hermes/workspace/repos/repociv/scripts/repociv-backup.sh
